# etc/systemd/user/bedtime.service
[Unit]
Description=Bedtime Enforcer Daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/bedtime
Restart=no
