# etc/systemd/user/free.service
[Unit]
Description=Time Guard

[Service]
ExecStart=%h/.local/bin/free guard

[Install]
WantedBy=default.target
